# cmd/f2k/test4.timer
[Unit]
Description=Test timer
Documentation=https://git.colofoo.net/fleet/test4

[Timer]
OnCalendar=*-*-* 6:00:00

[X-Fleet]
Global=true
